# infra/freeswitch/webhook-shim/voizo-shim.service
# Voizo FreeSWITCH → Dashboard webhook shim — systemd unit.
#
# Deploy to /etc/systemd/system/voizo-webhook-shim.service on the AWS EC2
# instance after copying the shim into /opt/voizo-webhook-shim.
#
# Deployment steps (run on the EC2 box):
#   sudo mkdir -p /opt/voizo-webhook-shim
#   sudo cp -r index.js package.json /opt/voizo-webhook-shim/
#   cd /opt/voizo-webhook-shim && sudo npm install --production
#   sudo cp voizo-shim.service /etc/systemd/system/voizo-webhook-shim.service
#   sudo nano /etc/systemd/system/voizo-webhook-shim.service  # fill in env vars
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now voizo-webhook-shim
#   sudo journalctl -u voizo-webhook-shim -f  # watch logs

[Unit]
Description=Voizo FreeSWITCH webhook shim — translates FS ESL events to Voizo dashboard webhooks
After=network.target freeswitch.service
Requires=freeswitch.service

[Service]
Type=simple
WorkingDirectory=/opt/voizo-webhook-shim
ExecStart=/usr/bin/node /opt/voizo-webhook-shim/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Runs as the EC2 default 'admin' user (Debian-based AMI), matching the
# originate-shim service. Non-root for defense-in-depth.
User=admin
Group=admin

# Environment — fill these in with real values before enabling the service.
Environment="ESL_HOST=127.0.0.1"
Environment="ESL_PORT=8021"
Environment="ESL_PASSWORD=REPLACE_ME_match_event_socket.conf.xml"
Environment="VOIZO_WEBHOOK_URL=https://voizo-eight.vercel.app/api/webhooks/freeswitch/voice-status"
Environment="VOIZO_WEBHOOK_SECRET=REPLACE_ME_random_hex_match_Vercel_FREESWITCH_WEBHOOK_SECRET"

[Install]
WantedBy=multi-user.target
